# Kanban/deployment/kanban.service
[Unit]
Description=Kanban Spring Boot Service
After=network.target docker.service

[Service]
User=your-user
WorkingDirectory=/home/your-user/kanban-project/deploy
# Load your secrets into the process environment
EnvironmentFile=/home/your-user/kanban-project/deploy/.env
ExecStart=/usr/bin/java -jar /home/your-user/kanban-project/target/kanban-app.jar
SuccessExitStatus=143
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
